# linux/Systemd/example.service
[Unit]
Description=Eample Service

[Service]
ExecStart=/usr/sbin/example

[Install]
WantedBy=multi-user.target